# systemd/manageiq-providers-vmware_cloud_manager_event_catcher@.service
[Unit]
PartOf=manageiq-providers-vmware_cloud_manager_event_catcher.target
[Install]
WantedBy=manageiq-providers-vmware_cloud_manager_event_catcher.target
[Service]
WorkingDirectory=/var/www/miq/vmdb
Environment=BUNDLER_GROUPS=manageiq_default,ui_dependencies
ExecStart=/bin/bash -lc 'exec ruby lib/workers/bin/run_single_worker.rb ManageIQ::Providers::Vmware::CloudManager::EventCatcher --heartbeat --guid=%i'
Restart=no
Type=notify
Slice=manageiq-providers-vmware_cloud_manager_event_catcher.slice
